# deploy/com.openbanking-mcp.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Plantilla de LaunchAgent (macOS). Sustituye __PROJECT_DIR__ por la ruta
  absoluta donde clonaste el repo, por ejemplo /Users/tu-usuario/openbanking-mcp
  Tip: sed -i '' "s|__PROJECT_DIR__|$PWD|g" deploy/com.openbanking-mcp.sync.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openbanking-mcp.sync</string>

    <!-- Ejecuta `finmcp sync` desde el venv del proyecto -->
    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/.venv/bin/finmcp</string>
        <string>sync</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- Cada 6 horas (21600 s). El refresh_token mantiene viva la sesión PSD2. -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <!-- Lanza una sincronización también al cargar el agente -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/data/sync.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/data/sync.err.log</string>
</dict>
</plist>
